# etc/cron.d/vzubcmon
# crontab file for vzubcmon

0 * * * *	root	python /usr/local/vzubcmon/vzubcmon

